# templates/borgmatic.timer.j2
[Unit]
Description=Run borgmatic backup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=60min

[Install]
WantedBy=timers.target
