# templates/borgmatic.timer.j2
[Unit]
Description=Run borgmatic backup

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=10min

[Install]
WantedBy=timers.target
